# crontab
* * * * * cd /app && php artisan schedule:work >> /dev/null 2>&1
